# rootfs/etc/services.d/homeseer/run
#!/usr/bin/with-contenv sh

# wait for services to come up
sleep 2

echo "Starting HomeSeer..."
cd /HomeSeer
exec mono ./HSConsole.exe --log
